# etc/supervisord.conf
[program:gunnery]
command = uwsgi --emperor /etc/uwsgi/apps-enabled/gunnery.ini

[program:worker]
directory = /var/gunnery/gunnery/
command = python manage.py celery worker
